# system_files/shared/usr/lib/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noreset --noclear --autologin ctv - ${TERM}
Environment=XDG_SESSION_TYPE=wayland
